# files/check-mediaserver.cron
*/10 * * * * /usr/lib/mediaserver/check-mediaserver
