# files/restart-mediaserver.cron
0 4 * * * systemctl restart mediaserver
